# bgalin-backend.service
[Unit]
Description=BGalin Backend Service (Bun/Elysia)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=bgalin
WorkingDirectory=/var/www/bgalin/server
# Try multiple Bun paths for flexibility
ExecStart=/bin/bash -c "exec /home/bgalin/.bun/bin/bun run src/index.ts"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="PORT=8000"
Environment="PATH=/home/bgalin/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/bgalin

[Install]
WantedBy=multi-user.target
