# bgalin-backend.service
[Unit]
Description=BGalin Backend Service (Bun/Elysia)
After=network.target

[Service]
Type=simple
User=bgalin
WorkingDirectory=/var/www/bgalin/server

# Bun executable path with full permissions
ExecStart=/bin/bash -c "exec /home/bgalin/.bun/bin/bun run src/index.ts"

Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal

Environment="NODE_ENV=production"
Environment="PORT=8000"
Environment="BUN_INSTALL=/home/bgalin/.bun"
Environment="PATH=/home/bgalin/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/bgalin"

# Limit memory
Environment="NODE_OPTIONS=--max-old-space-size=256"

[Install]
WantedBy=multi-user.target
